# net/openflow-ucap/files/etc/init.d/openflow
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org

START=98
STOP=02

boot() {
	ofup -a
}

start() {
	ofup -a
}

restart() {
	stop
	start
}

stop() {
	ofdown -a
}
