# utils/bismark-experiments-manager/files/etc/init.d/bismark-experiments-manager
#!/bin/sh /etc/rc.common

START=99
start() {
    [ -x /usr/bin/bismark-experiments-manager ] && \
        /usr/bin/bismark-experiments-manager
}

stop() {
    :
}
