# utils/bismark-experiments-manager/files/etc/init.d/bismark-experiments-manager
#!/bin/sh /etc/rc.common

START=98
start() {
    [ -x /usr/bin/merge-experiments ] && \
        /usr/bin/merge-experiments
}

stop() {
    :
}
