# utils/bismark-updater/files/etc/cron.d/bismark-updater
0 5,17 * * * /usr/bin/bismark-updater -d
