# utils/bismark-registration/files/etc/uci-defaults/bismark-authn
#!/bin/sh
/etc/init.d/bismark-authn-check enable
(crontab -l | grep -v bismark-authn-check; echo '0 0 * * * /usr/bin/bismark-authn-check -d') | crontab -
mv /etc/nodogsplash/htdocs/images/nodogsplash.conf /etc/nodogsplash/
mv /etc/nodogsplash/htdocs/images/splash.html /etc/nodogsplash/htdocs/
/etc/init.d/nodogsplash restart
cat /www/cgi-bin/projectbismark-verify | sed 's/myrouter/register/' | sed 's/http/https/' >"/www/cgi-bin/register-verify"
chmod a+x /www/cgi-bin/register-verify
