# utils/bismark-registration/files/etc/uci-defaults/bismark-authn
#!/bin/sh

# Enable init script
/etc/init.d/bismark-authn-check enable

# Enable cron job
crontab -l | cat - /etc/crontabs.d/bismark-registration | sort -u | crontab -

# Sort files, start nodogsplash
mv /etc/nodogsplash/nodogsplash.conf.will-rename-in-uci-defaults /etc/nodogsplash/nodogsplash.conf
mv /etc/nodogsplash/htdocs/splash.html.will-rename-in-uci-defaults /etc/nodogsplash/htdocs/splash.html
/etc/init.d/nodogsplash restart
