# utils/bismark-registration/files/etc/init.d/bismark-authn-check
#!/bin/sh /etc/rc.common

#checks for authentication token on startup

START=98

boot() {
    bismark-authn-check
}
