# utils/bismark-registration/files/etc/init.d/bismark-authn-check
#!/bin/sh /etc/rc.common

START=66
STOP=70

start() {
	. /usr/bin/bismark-authn-check
}
