# utils/bismark-updater/files/etc/uci-defaults/bismark-updater
#!/bin/sh
/etc/init.d/bismark-updater enable
(crontab -l | grep -v bismark-updater; echo '0 5,17 * * * /usr/bin/bismark-updater -d') | crontab -
