# utils/bismark-updater/files/etc/uci-defaults/bismark-updater
#!/bin/sh
/etc/init.d/bismark-updater enable
crontab -l | cat - /etc/crontabs.d/bismark-updater | sort -u | crontab -
