# utils/bismark-health/files/etc/uci-defaults/bismark-health
#!/bin/sh
(crontab -l | grep -v bismark-health; echo '0 * * * * /usr/bin/bismark-health -d') | crontab -
