# utils/bismark-health/files/etc/uci-defaults/bismark-health
#!/bin/sh
crontab -l | cat - /etc/crontabs.d/bismark-health | sort -u | crontab -
