# utils/bismark-censorship-performance-measurements-tmpfs/files/etc/cron.d/censorship-performance-measurements
*/10 * * * * /tmp/usr/bin/censorship-performance-measurements-http
