# utils/bismark-experiments-manager/files/etc/uci-defaults/bismark-experiments-manager
#!/bin/sh
/etc/init.d/bismark-experiments-manager enable || true
