# utils/bismark-experiments-manager/files/etc/uci-defaults/bismark-experiments-manager
#!/bin/sh
/etc/init.d/bismark-experiments-manager enable
(crontab -l | grep -v merge-experiments; echo '30 * * * * /usr/bin/merge-experiments') | crontab -
