# utils/bismark-experiments-manager/files/etc/uci-defaults/bismark-experiments-manager
#!/bin/sh
/etc/init.d/bismark-experiments-manager enable
crontab -l | cat - /etc/crontabs.d/bismark-experiments-manager | sort -u | crontab -
