# utils/bismark-experiments-manager/files/etc/uci-defaults/bismark-experiments-manager
#!/bin/sh
/etc/init.d/bismark-experiments-manager enable
(crontab -l | grep -v bismark-experiments-manager; echo '0 5,17 * * * /usr/bin/bismark-experiments-manager -d') | crontab -
