# utils/pakistan-censorship-testing/files/etc/pakistan-censorship-testing.crontab
*/10 3,9,15,21 * * * /tmp/usr/bin/pakistan-censorship-testing
