# utils/censorship-performance/files/etc/censorship-performance.crontab
*/10 * * * /tmp/usr/bin/censorship-performance-http
